# app/layout/DEBIAN/postinst
#!/bin/sh
set -e

PATH="/var/jb/usr/bin:/var/jb/bin:/usr/bin:/bin:/usr/sbin:/sbin"
export PATH

BIN="/var/jb/usr/libexec/ipadecryptd"
SUPERVISOR="/var/jb/usr/libexec/ipadecryptd-supervisor"
PLIST="/var/jb/Library/LaunchDaemons/com.korboy.ipadecryptd.plist"
OUTDIR="/var/mobile/Documents/ipadecrypt/decrypted"

launchctl bootout system "$PLIST" 2>/dev/null || true
ps -axo pid,command | awk '/[i]padecryptd-supervisor/ { print $1 }' | while read -r pid; do
    kill "$pid" 2>/dev/null || true
done
killall ipadecryptd 2>/dev/null || true
rm -f /var/jb/var/run/ipadecryptd.sock 2>/dev/null || true

chmod 755 "$BIN" "$SUPERVISOR"
chown 0:0 "$BIN" "$SUPERVISOR"
chmod 644 "$PLIST"
chown 0:0 "$PLIST"
mkdir -p "$OUTDIR"
chown mobile:mobile /var/mobile/Documents/ipadecrypt "$OUTDIR"
chmod 755 "$OUTDIR"

launchctl bootstrap system "$PLIST" 2>/dev/null || launchctl load "$PLIST" 2>/dev/null || true
launchctl kickstart -k system/com.korboy.ipadecryptd 2>/dev/null || true

exit 0
